# api/systemd/rpo-v5-api.service
[Unit]
Description=RPO V5 API - PostgreSQL REST API (porta 7001)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cazouvilela
WorkingDirectory=/home/cazouvilela/projetos/RPO/api
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
