# api/systemd/rpo-v5-worker.service
[Unit]
Description=RPO V5 Worker - Valkey to PostgreSQL Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cazouvilela
WorkingDirectory=/home/cazouvilela/projetos/RPO/api
ExecStart=/usr/bin/node worker.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
